# arm/rootfs-arm/etc/inittab
::sysinit:/etc/init.d/rcS
#::askfirst:-/bin/ash
::respawn:/sbin/getty 115200 console
::respawn:/sbin/syslogd -n
